# dot_config/systemd/user/uv-cache-prune.service
[Unit]
Description=uv cache size prune

[Service]
Type=oneshot
ExecStart=%h/.local/bin/uv-cache-check
