# dot_config/autostart/fastfetch.desktop
[Desktop Entry]
Type=Application
Name=Fastfetch
Exec=wezterm start -- bash -c "fastfetch; exec zsh"
X-KDE-autostart-phase=2
