# dot_config/autostart/fastfetch.desktop
[Desktop Entry]
Type=Application
Name=Fastfetch
Exec=~/.config/autostart/fastfetch-launch.sh
X-KDE-autostart-phase=2
